# network/inetutils/files/init.d/rlogind
#!/sbin/openrc-run

command=-/usr/bin/rlogind
command_args=""
pidfile=

name="rlogind"
description="Remote Login Facilities Server"
