# network/inetutils/files/init.d/telnet
#!/sbin/openrc-run

command=-/usr/bin/telnetd
command_args=""
pidfile=

name="telnetd"
description="Telnet Server"
